# api.socket
[Unit]
Description=api sock
[Socket]
ListenStream=/run/api.sock
[Install]
WantedBy=sockets.target