# fixie.socket
[Unit]
Description=fixie socket
[Socket]
ListenStream=/run/fixie.sock
[Install]
WantedBy=sockets.target