# fixie.socket
[Unit]
Description=fixie sock
[Socket]
ListenStream=/run/fixie.sock
[Install]
WantedBy=sockets.target